# deployment/django.service
[Unit]
Description=Smart Camera Django Application (Gunicorn)
After=network.target postgresql.service redis-server.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/smartcamera/django-webapp
Environment="PATH=/var/www/smartcamera/django-webapp/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=iot_platform.settings_production"

# Load environment variables from .env file
EnvironmentFile=/var/www/smartcamera/django-webapp/.env

# Gunicorn executable
ExecStart=/var/www/smartcamera/django-webapp/venv/bin/gunicorn \
    --workers 3 \
    --worker-class sync \
    --timeout 120 \
    --bind 127.0.0.1:8000 \
    --access-logfile /var/log/smartcamera/gunicorn-access.log \
    --error-logfile /var/log/smartcamera/gunicorn-error.log \
    --log-level info \
    iot_platform.wsgi:application

# Restart on failure
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target







